# linux-cortexm/local/etc/inittab
::sysinit:/etc/rc
ttyS0::respawn:/bin/hush -i
